# LAB2/dist/rootfs/etc/inittab
ttyS0::respawn:/bin/getty -L ttyS0 115200
